# files/default/system.conf
respawn

start on runlevel [2345]
stop on runlevel [016]

exec system --statsd-address localhost:8125
